# docs/wyrmlock.service
[Unit]
Description=wyrmlock Application Security Daemon
Documentation=https://github.com/ZanzyTHEbar/wyrmlock
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/etc/wyrmlock
ExecStart=/usr/local/bin/wyrmlock daemon
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5
# Give the service 3 seconds to stop gracefully before killing it
KillSignal=SIGINT
TimeoutStopSec=3

# Security hardening
CapabilityBoundingSet=CAP_SYS_PTRACE CAP_KILL CAP_SYS_ADMIN
AmbientCapabilities=CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_SYS_ADMIN
SecureBits=keep-caps
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/var/run /etc/wyrmlock
ProtectHome=read-only
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_UNIX AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
MemoryDenyWriteExecute=true
PrivateTmp=true
RestrictSUIDSGID=true
LockPersonality=true

# Limit resource usage
LimitNOFILE=1024
MemoryLimit=256M

[Install]
WantedBy=multi-user.target
